# launchd/com.elijah.tradingagent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>com.elijah.tradingagent</string>
  <key>ProgramArguments</key><array><string>/Users/elijahang/TradingAgentRuntime/scripts/run_once.sh</string></array>
  <key>WorkingDirectory</key><string>/Users/elijahang/TradingAgentRuntime</string>
  <key>StartInterval</key><integer>600</integer>
  <key>StandardOutPath</key><string>/Users/elijahang/Library/Application Support/TradingAgent/logs/launchd.out</string>
  <key>StandardErrorPath</key><string>/Users/elijahang/Library/Application Support/TradingAgent/logs/launchd.err</string>
</dict></plist>
