# launchd/com.elijah.tradingagent.telegram.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0"><dict>
  <key>Label</key><string>com.elijah.tradingagent.telegram</string>
  <key>ProgramArguments</key><array><string>/Users/elijahang/TradingAgentRuntime/scripts/run_telegram_listener.sh</string></array>
  <key>WorkingDirectory</key><string>/Users/elijahang/TradingAgentRuntime</string>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>/Users/elijahang/Library/Application Support/TradingAgent/logs/listener-launchd.out</string>
  <key>StandardErrorPath</key><string>/Users/elijahang/Library/Application Support/TradingAgent/logs/listener-launchd.err</string>
</dict></plist>
